# dist-install-templates/lin64/prime/daemon/install.sh
sudo cp prime.service /etc/systemd/system/
sudo systemctl enable prime.service
sudo systemctl status prime.service
